# bin/bs2b/bs2b@.service
[Unit]
Description=bs2b
Requires=pulseaudio.service
After=pulseaudio.service

[Service]
Type=oneshot
ExecStart=/home/daniel/bin/bs2b %i start
ExecStop=/home/daniel/bin/bs2b %i stop
RemainAfterExit=yes

[Install]
WantedBy=default.target
